# downloads-organizer.service.template
[Unit]
Description=Hourly Downloads Folder Organizer
After=network.target

[Service]
Type=oneshot
ExecStart=/home/c7/.scripts/downloads_organizer.sh

[Install]
WantedBy=default.target
